# telegram-dl-bot.service
[Unit]
Description=Telegram Video Downloader Bot
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /path/to/telegram_video_downloader.py
WorkingDirectory=/path/to/
EnvironmentFile=/path/to/.env
Restart=always
User=your_user
Group=your_user

[Install]
WantedBy=multi-user.target
